# conf/php_supervisord.conf
[program:php-fastcgi]
command=php-cgi -b /tmp/php-fastcgi.sock
directory=/
user=www-data
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/var/log/php.log